# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
loglevel=info

; ── PHP-FPM ──────────────────────────────────────────────────────────────────
[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
priority=10
startretries=5
startsecs=2
stderr_logfile=/var/log/supervisor/php-fpm.err.log
stdout_logfile=/var/log/supervisor/php-fpm.out.log
stdout_logfile_maxbytes=5MB

; ── Nginx ─────────────────────────────────────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
startretries=5
startsecs=2
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
stdout_logfile_maxbytes=5MB

; ── Horizon is intentionally removed from this supervisor config ───────────────
; Run Horizon as a SEPARATE EasyPanel service using the same Docker image
; with the start command:  php /var/www/html/artisan horizon
; This prevents horizon crashes from taking down nginx + php-fpm.
