# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
loglevel=info

; ── PHP-FPM ──────────────────────────────────────────────────────────────────
[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
priority=10
startretries=5
startsecs=2
stderr_logfile=/var/log/supervisor/php-fpm.err.log
stdout_logfile=/var/log/supervisor/php-fpm.out.log
stdout_logfile_maxbytes=5MB

; ── Nginx ─────────────────────────────────────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
startretries=5
startsecs=2
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
stdout_logfile_maxbytes=5MB

; ── Queue Worker ──────────────────────────────────────────────────────────────
; Processes background jobs (AnalyzeConversation runs on the "analyses" queue).
; Uses plain queue:work (not Horizon) so it stays simple and self-contained.
; autorestart=true keeps it alive; if Redis is briefly unavailable it retries.
[program:queue-worker]
command=php /var/www/html/artisan queue:work redis --queue=analyses,default --sleep=3 --tries=3 --max-time=3600 --timeout=310
autostart=true
autorestart=true
priority=30
startsecs=3
startretries=10
user=www-data
stopwaitsecs=330
killasgroup=true
stderr_logfile=/var/log/supervisor/queue.err.log
stdout_logfile=/var/log/supervisor/queue.out.log
stdout_logfile_maxbytes=10MB
